# root-fs/etc/local.d/ksm.start
#!/bin/bash
echo 1 > /sys/kernel/mm/ksm/run ;
exit
